# systemd/avologo-client.service
[Unit]
Description=Avologo client service

[Service]
Type=simple
ExecStart=/opt/avologo/avologo --mode client

[Install]
WantedBy=multi-user.target